# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:postgres]
command=su - postgres -c "/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/main -c config_file=/etc/postgresql/15/main/postgresql.conf"
priority=10
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.err

[program:redis]
command=/usr/bin/redis-server --bind 0.0.0.0 --protected-mode no
priority=20
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err

[program:backend]
command=/usr/local/bin/uv run python manage.py runserver 0.0.0.0:%(ENV_API_PORT)s
directory=/app/backend
environment=PYTHONPATH="/app/backend",DATABASE_URL="%(ENV_DATABASE_URL)s"
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.err

[program:frontend]
command=node server.js
directory=/app/frontend
environment=PORT="%(ENV_PORT)s",HOSTNAME="0.0.0.0"
priority=40
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend.err
